# systemd/ulanzi-daemon.service
[Unit]
Description=Ulanzi D200 Manager Daemon
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/ulanzi-daemon %h/.config/ulanzi/config.yaml
Restart=on-failure
RestartSec=5

# Device access
DevicePolicy=auto
DeviceAllow=char-hidraw rw

[Install]
WantedBy=default.target
